# profile.d/images.sh
####################################### Autogenerated by: images.sh
# IMAGES:   Images Names (IMAGE_NAME -> DOCKER_REPO:DOCKER_TAG)
export IMAGES
IMAGES=$(cat <<eot
j5pu/rc:alpine
eot
)
