# profile.d/tags.sh
TAGS=$(cat <<EOF
j5pu/alpine:latest
j5pu/archlinux:latest
j5pu/bash:latest
j5pu/bats:latest
j5pu/busybox:latest
j5pu/centos:latest
j5pu/debian:latest
j5pu/fedora:latest
j5pu/systemd-ubuntu:latest
j5pu/kali-rolling:latest
j5pu/kali-bleeding-edge:latest
j5pu/nix:latest
j5pu/python:3.9-alpine
j5pu/python:3.9-bullseye
j5pu/python:3.9-slim
j5pu/python:3.10-alpine
j5pu/python:3.10-bullseye
j5pu/python:3.10-slim
j5pu/rhel7:latest
j5pu/ubuntu:latest
j5pu/zsh:latest
EOF
)
export TAGS
